# pi-purify.service
[Unit]
Description=pi-purify

[Service]
ExecStart=/usr/local/bin/pi-purify.py

[Install]
WantedBy=multi-user.target
